# contrib/packaging/debian/clacks/debian/supervisor/clacks-client.conf
[program:clacks-client]
command=/usr/sbin/clacks-client
user=clacks
stopsignal=INT
