# contrib/packaging/debian/clacks/debian/supervisor/clacks-dbus.conf
[program:clacks-dbus]
command=/usr/sbin/clacks-dbus
stopsignal=INT
