# contrib/packaging/debian/clacks/debian/supervisor/clacks-agent.conf
[program:clacks-agent]
command=/usr/sbin/clacks-agent
user=clacks
stopsignal=INT
